# conf/supervisor/tornado-demo.conf
[program:app]
command=python /var/www/tornado-demo/web.py --port=8001
directory=/var/www/tornado-demo
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/tornado-demo.log
